# service/buildagent/run
#!/bin/bash
set -e
echo "Starting buildagent..."
chown -R teamcity:root /opt/buildAgent
chown -R teamcity:root /data
exec /sbin/setuser teamcity /opt/buildAgent/bin/agent.sh run
